# utils/wccf_proc_multiap.in
#!/bin/sh /etc/rc.common
#
# Copyright (c) 2019 Cable Television Laboratories, Inc.
# Licensed under the BSD+Patent (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#          https://opensource.org/licenses/BSDplusPatent
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
# Init script: wccf_proc_multiap
# usage: $ /etc/init.d/wccf_proc_multiap [ enable | disable | start | stop | restart ]
# To be placed in /etc/init.d

START=90
STOP=10

source /etc/config/wccf_proc_multiap.cfg

# From above cfg file (first two may differ):
# WCCF_PIDFILE=/tmp/run/wccf_proc_multiap.pid
# WCCF_INPUT_PATH=/tmp/wccf/sensor
# WCCF_OUTPUT_PATH=/tmp/wccf/tx
# WCCF_LOG_PATH=/root/log


start() {
	if [ -f "$WCCF_PIDFILE" ]; then
		echo "File $WCCF_PIDFILE exists - cannot start wccf_proc_multiap"
		exit 1
	fi
	echo "Starting wccf_proc_multiap"
	mkdir -p $WCCF_INPUT_PATH
	mkdir -p $WCCF_OUTPUT_PATH
	mkdir -p $WCCF_LOG_PATH
	/usr/sbin/wccf_proc_multiap --input-path=$WCCF_INPUT_PATH --output-path=$WCCF_OUTPUT_PATH >> $WCCF_LOG_PATH/wccf_proc_multiap.log 2>&1 & echo $! > $WCCF_PIDFILE
	PID=`cat $WCCF_PIDFILE`
	echo "Process wccf_proc_multiap started with pid $PID"
}

stop() {
	if [ ! -f "$WCCF_PIDFILE" ]; then
		echo "File $WCCF_PIDFILE missing - cannot stop wccf_proc_multiap"
		exit 1
	fi
	echo "Stoping wccf_proc_multiap"
	WCCF_PID=`cat $WCCF_PIDFILE`
	kill -9 $WCCF_PID
	rm -f $WCCF_PIDFILE
}

restart() {
	stop
	start
}

