# utils/wccf_tx_rest.in
#!/bin/sh /etc/rc.common
#
# Copyright (c) 2019 Cable Television Laboratories, Inc.
# Licensed under the BSD+Patent (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#          https://opensource.org/licenses/BSDplusPatent
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
# Init script: wccf_tx_rest
# usage: $ /etc/init.d/wccf_tx_rest [ enable | disable | start | stop | restart ]
# To be placed in /etc/init.d

START=90
STOP=10

source /etc/config/wccf_tx_rest.cfg

# From above cfg file (first two may differ):
# WCCF_RECEIVER_HOST=208.113.103.113
# WCCF_RECEIVER_PORT=8888
# WCCF_PIDFILE=/tmp/run/wccf_tx_rest.pid
# WCCF_INPUT_PATH=/tmp/wccf/tx/rest
# WCCF_LOG_PATH=/root/log


start() {
	if [ -f "$WCCF_PIDFILE" ]; then
		echo "File $WCCF_PIDFILE exists - cannot start wccf_tx_rest"
		exit 1
	fi
	echo "Starting wccf_tx_rest"
	mkdir -p $WCCF_INPUT_PATH
	mkdir -p $WCCF_LOG_PATH
	/usr/sbin/wccf_tx_rest --input-path=$WCCF_INPUT_PATH --host=$WCCF_RECEIVER_HOST --port=$WCCF_RECEIVER_PORT >> $WCCF_LOG_PATH/wccf_tx_rest.log 2>&1 & echo $! > $WCCF_PIDFILE
	PID=`cat $WCCF_PIDFILE`
	echo "Process wccf_tx_rest started with pid $PID"
}

stop() {
	if [ ! -f "$WCCF_PIDFILE" ]; then
		echo "File $WCCF_PIDFILE missing - cannot stop wccf_tx_rest"
		exit 1
	fi
	echo "Stoping wccf_tx_rest"
	WCCF_PID=`cat $WCCF_PIDFILE`
	kill -9 $WCCF_PID
	rm -f $WCCF_PIDFILE
}

restart() {
	stop
	start
}

