# utils/wccf_jobs.cron.in
#
# Replace these times with appropriate settings for your use case.
# Set output-dir as appropriate.
#
*/5 * * * * /usr/sbin/wccf_sensor_scan --interface-name=wlan0 --output-dir=/tmp/wccf/sensor
*/5 * * * * /usr/sbin/wccf_sensor_scan --interface-name=wlan1 --output-dir=/tmp/wccf/sensor

*/5 * * * * /usr/sbin/wccf_sensor_survey --interface-name=wlan0 --output-dir=/tmp/wccf/sensor
*/5 * * * * /usr/sbin/wccf_sensor_survey --interface-name=wlan1 --output-dir=/tmp/wccf/sensor

*/5 * * * * /usr/sbin/wccf_sensor_station --interface-name=wlan0 --output-dir=/tmp/wccf/sensor
*/5 * * * * /usr/sbin/wccf_sensor_station --interface-name=wlan1 --output-dir=/tmp/wccf/sensor

# File count and target directory based on above settings.
#
*/10 * * * * /usr/sbin/wccf_util_purge.sh 50 /tmp/wccf/sensor

